# templates/etc/init/riemann-dash.conf.erb
description "start and stop riemann-dash"
author "Gareth Rushgrove"

start on runlevel [2345]
stop on runlevel [!2345]

setuid <%= user %>

respawn
respawn limit 5 20

exec riemann-dash /etc/riemann-dash.rb
